# de.agorion.tcrb.sample.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Unique label in reverse-domain notation -->
    <key>Label</key>
    <string>de.agorion.tcrb</string>

    <!-- What gets executed. Adjust paths!
         check python3 with:  which python3   (Apple system Python: /usr/bin/python3) -->
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>/Users/USERNAME/Scripts/tcrb/tcrb_monitor.py</string>
    </array>

    <!-- Working directory: tcrb_history.csv and tcrb_state.json are written here -->
    <key>WorkingDirectory</key>
    <string>/Users/USERNAME/Scripts/tcrb</string>

    <!-- Hourly on the hour. launchd catches up a missed run
         as soon as the Mac wakes from sleep. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- Run once immediately on load (useful for testing; optional later). -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Logs for troubleshooting -->
    <key>StandardOutPath</key>
    <string>/Users/USERNAME/Scripts/tcrb/tcrb.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/USERNAME/Scripts/tcrb/tcrb.err.log</string>

    <!-- launchd provides only a minimal PATH. osascript (for the
         macOS notification) lives in /usr/bin. -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
